# security/pfSense-pkg-crowdsec/pkg-plist
pkg/crowdsec.xml
pkg/crowdsec.inc
etc/rc.conf.d/pfsense_crowdsec
etc/rc.d/pfsense_crowdsec
%%DATADIR%%/info.xml
