# sysutils/pfSense-azure/pkg-plist
/etc/waagent.conf
etc/rc.d/waagent.sh
sbin/add-pfsense-user
sbin/set-pfsense-hostname
sbin/set-pfsense-password
sbin/set-pfsense-sshkey
sbin/waagent
@dir /etc
